# templates/pg_backup.cron.j2
# min hour daymo mo daywk user cmd
1 2 * * * {{ db.backups.user }} {{ db.backups.rootdir }}/pg_backup.sh
