# files/solr.conf
start on runlevel [2345]
stop on runlevel [!2345]

kill timeout 30
respawn

#setuid solr
#setgid solr

script
   chdir /usr/local/solr
   exec /usr/bin/java -jar start.jar
end script
